# docker/supervisord.conf
[supervisord]
nodaemon=true
[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -noxrecord -noxfixes -noxdamage -forever -rfbauth ~/.vnc/passwd
